# sandbox-image/etc/profile.d/common-bashrc.sh
if [ -f /shared/common/bashrc ]; then
  . /shared/common/bashrc
fi
